# printer-api-restart.service
[Unit]
Description=Restart Printer API Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart printer-api.service